# bin/cluster
#!/bin/bash
### BEGIN INIT INFO
# Provides:          thin
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: cluster initscript
# Description:       cluster
### END INIT INFO

# Original author: Forrest Robertson
# Modified by: Hampton Catlin

trap "echo Trapped!" USR1

# Do NOT "set -e"

DAEMON=/srv/mobile/current/bin/server
SCRIPT_NAME=/srv/mobile/current/bin/cluster
CONFIG_PATH=/srv/mobile/current/config/thins

# Exit if the package is not installed
[ -x "$DAEMON" ] || exit 0

case "$1" in
  start)
	$DAEMON start --all $CONFIG_PATH
	;;
  stop)
	pkill -9 thin
	;;
  restart)
	$DAEMON restart -O --all $CONFIG_PATH
	;;
  *)
	echo "Usage: $SCRIPT_NAME {start|stop|restart}" >&2
	exit 3
	;;
esac

:
